# lxc-automat.service
[Unit]
Description=LXC-Automat Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lxc-automat
ExecStart=/opt/lxc-automat/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8091
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
